# deployment/systemd/dayplanner-backend.service
[Unit]
Description=DayPlanner Backend API Service
After=network.target mysql.service
Wants=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/dayplanner/backend
Environment=PYTHONPATH=/opt/dayplanner/backend
EnvironmentFile=/opt/dayplanner/.env.production
ExecStartPre=/opt/dayplanner/deployment/scripts/prestart-backend.sh
ExecStart=/opt/dayplanner/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dayplanner-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dayplanner /var/log/dayplanner /tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
